# bridge/com.molant.ha-lucarne-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.molant.ha-lucarne-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/shortcuts</string>
    <string>run</string>
    <string>ha-lucarne-sync</string>
  </array>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/<USERNAME>/Library/Logs/ha-lucarne-sync.log</string>
  <key>StandardErrorPath</key><string>/Users/<USERNAME>/Library/Logs/ha-lucarne-sync.log</string>
</dict>
</plist>
